# deploy/logclean.service
[Unit]
Description=Fiesta tracker log rotation (keep < 1 month, total logs <= 10 GB)

[Service]
Type=oneshot
ExecStart=/root/logclean.sh
